# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=su - postgres -c "/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf"
autostart=true
autorestart=true
priority=1
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server --protected-mode no
autostart=true
autorestart=true
priority=2
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth]
command=python main.py
directory=/app/user_auth
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=10
environment=DATABASE_USER="postgres",DATABASE_PASSWORD="asknova123",DATABASE_HOST="localhost",DATABASE_PORT="5432"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sessions]
command=python main.py
directory=/app/session_chat
autostart=true
autorestart=true
priority=11
startsecs=5
startretries=10
environment=DATABASE_USER="postgres",DATABASE_PASSWORD="asknova123",DATABASE_HOST="localhost",DATABASE_PORT="5432",USER_SERVICE_URL="http://localhost:8006",REDIS_HOST="localhost"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:messages]
command=python main.py
directory=/app/messages
autostart=true
autorestart=true
priority=12
startsecs=5
startretries=10
environment=DATABASE_USER="postgres",DATABASE_PASSWORD="asknova123",DATABASE_HOST="localhost",DATABASE_PORT="5432",SESSION_SERVICE_URL="http://localhost:8005"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stt]
command=python main.py
directory=/app/stt
autostart=true
autorestart=true
priority=20
startsecs=10
startretries=5
environment=STT_HOST="0.0.0.0",STT_PORT="8001",DEVICE="cuda"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cag]
command=python main.py
directory=/app/cag
autostart=true
autorestart=true
priority=21
startsecs=10
startretries=5
environment=PORT="8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=python -m uvicorn gateway.gateway:app --host 0.0.0.0 --port 8090 --workers 1
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=5
startretries=10
environment=STT_WS_URL="ws://localhost:8001/stream/mux",CAG_WS_URL="ws://localhost:8000/chat/ws",USER_SERVICE_URL="http://localhost:8006",SESSION_SERVICE_URL="http://localhost:8005",MESSAGE_SERVICE_URL="http://localhost:8003"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=npx tsx server.ts
directory=/app/frontend
autostart=true
autorestart=true
priority=40
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
